# apps/hub/deploy/zzapi-mes-hub-backup.timer
[Unit]
Description=Daily SQLite backup for ZZAPI MES Hub

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
