# apps/hub/deploy/zzapi-mes-hub-backup.service
[Unit]
Description=ZZAPI MES Hub SQLite backup
After=zzapi-mes-hub.service

[Service]
Type=oneshot
User=zzapi-mes
Group=zzapi-mes
TimeoutStartSec=300
ExecStart=/usr/local/bin/zzapi-mes-hub-backup.sh
